# src/backend/dbus/proxyd.service.in
[Unit]
Description=Libproxy user service

[Service]
Type=dbus
BusName=org.libproxy.proxy
ExecStart=@LIBEXECDIR@/proxyd
#Environment="G_MESSAGES_DEBUG=all"
Restart=on-failure

MemoryDenyWriteExecute=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=full
RestrictRealtime=yes

[Install]
WantedBy=default.target